# Yocto/InitMngr/Systemd/helloworld.service
[Unit]
Description=Hello World startup

[Service]
ExecStart=/usr/bin/helloworld

[Install]
WantedBy=multi-user.target
